# deploy/activity-daily.service
[Unit]
Description=Activity Daily Spring Boot service
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/activity-daily/server-spring
EnvironmentFile=/opt/activity-daily/.env
ExecStart=/usr/bin/java -jar /opt/activity-daily/server-spring/activity-daily-server.jar
Restart=always
RestartSec=5
User=activitydaily
Group=activitydaily

[Install]
WantedBy=multi-user.target